# meta-fsl-imx6/recipes-bsp/advantech-emmc-download/advantech-emmc-download/emmc-download.service
[Unit]
Description=EMMC Download Service

[Service]
Type=oneshot
ExecStart=/emmc/emmc-download.sh

[Install]
WantedBy=multi-user.target
