# meta-fsl-imx6/recipes-bsp/advantech-initscripts/advantech-initscripts/ota-rootfs-update.service
[Unit]
Description=update ota rootfs service
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/lib/ota-rootfs-update.sh

[Install]
WantedBy=custom.target
