# meta-fsl-imx6/recipes-bsp/update-tools/update-tools-1.0/files/updateserver.service
[Unit]
Description=Advantech OTA Download Service 2

[Service]
Type=oneshot
ExecStart=/usr/bin/UpdateServer.py -k /etc/updateServerPublicKey.pem

[Install]
WantedBy=multi-user.target
